# lib/systemd/system/astral-weekly-report.service
[Unit]
Description=Astral Matrix Weekly Health Report

[Service]
Type=oneshot
ExecStart=/usr/local/bin/astral-report
