# astral-matrix_2.2.0_all/DEBIAN/prerm
#!/usr/bin/env bash
set -e
case "${1:-}" in
  remove|purge)
    systemctl disable --now astral-weekly-report.timer 2>/dev/null || true
    systemctl daemon-reload || true
    ;;
esac
exit 0
